# deploy/systemd/arachne-share.service
[Unit]
Description=Arachne public snapshot origin
After=network-online.target

[Service]
Type=simple
Environment=ARACHNE_DEPLOY_ENV=%h/.config/arachne/deployment.env
WorkingDirectory=%h/arachne
ExecStart=%h/arachne/bin/run-configured-service.sh share
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
